# social_game_typology_quiz.service
[Unit]
Description=Social Game Typology Quiz
Requires=docker.service
After=docker.service

[Service]
Type=simple
WorkingDirectory=/home/warelock/projects/social_game_typology
ExecStart=/usr/bin/docker compose up --build
ExecStop=/usr/bin/docker compose down
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target

